# Docker/Hadoop/scripts/start-hdfs.service
[Unit]
After=network.service

[Service]
ExecStart=/root/start-hdfs.sh

[Install]
WantedBy=default.target